# deploy/niftystrategist-ordernode.service
[Unit]
Description=NiftyStrategist Order Node (local)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/niftystrategist/backend
EnvironmentFile=/opt/niftystrategist/backend/.env
ExecStart=/opt/niftystrategist/backend/venv/bin/uvicorn order_node.app:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
