# deploy/niftystrategist.service
[Unit]
Description=NiftyStrategist
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/niftystrategist/backend
EnvironmentFile=/opt/niftystrategist/backend/.env
ExecStart=/opt/niftystrategist/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
